# notification_service/launchd/com.betclub.notifier.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.betclub.notifier</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/ryan.sites/Downloads/betting-app-dashboard/notification_service/build/macos/Build/Products/Release/notification_service.app/Contents/MacOS/notification_service</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<false/>
	<key>WorkingDirectory</key>
	<string>/Users/ryan.sites/Downloads/betting-app-dashboard/notification_service</string>
	<key>StandardOutPath</key>
	<string>/tmp/com.betclub.notifier.stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/com.betclub.notifier.stderr.log</string>
</dict>
</plist>
